# etc/systemd/system/systemd-nspawn@greenhouses.service.d/resources.conf
[Service]
CPUQuota=40%
MemoryMax=2G
MemoryHigh=1G
